# scripts/com.nitter.scrape.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nitter.scrape</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/varun/Projects/nitter/scripts/scrape_catchup.sh</string>
    </array>

    <!-- Run when user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Also try to run daily at 9am (catches laptop-open scenarios) -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Low priority - don't interfere with other tasks -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Logs for debugging -->
    <key>StandardOutPath</key>
    <string>/Users/varun/Projects/nitter/output/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/varun/Projects/nitter/output/launchd_stderr.log</string>
</dict>
</plist>
